# speech-api/supervisord.conf
[supervisord]
nodaemon=true
logfile=/workspace/speech-api/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
pidfile=/workspace/speech-api/supervisord.pid

[program:speech-api]
command=/workspace/speech-api/run.sh
directory=/workspace/speech-api
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
stdout_logfile=/workspace/speech-api/logs/speech-api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_redirect=true
environment=HOME="/root",HF_HOME="/tmp/hf_home",PATH="/workspace/speech-api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[supervisorctl]
serverurl=unix:///workspace/speech-api/supervisor.sock

[unix_http_server]
file=/workspace/speech-api/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
